# profile.d/alias
alias ll="ls -lGa"
alias la="ls -Ga"
alias l="ls -Ga"

alias h=history
alias hm="history | more"
alias tarc="tar -zcvf"
alias tarx="tar -zxvf"

alias cat=lolcat
alias start-proxy=cntlm

alias j=jobs
alias gr='"$(git rev-parse --show-toplevel)"/gradlew'
alias weather=wego
alias weather=wego
